# ja_uniforms.service
[Unit]
Description=J.A. Uniforms Pricing Application
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=notify
User=ubuntu
Group=www-data
WorkingDirectory=/workspace
Environment="PATH=/workspace/.venv/bin"
ExecStart=/workspace/.venv/bin/gunicorn --config gunicorn_config.py wsgi:app
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/workspace /var/log/ja_uniforms /var/run/ja_uniforms
PrivateDevices=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ja_uniforms

[Install]
WantedBy=multi-user.target
